# src/init.d/pixiv-backup
#!/bin/sh /etc/rc.common

START=95
STOP=01

USE_PROCD=1
PROG=/usr/bin/pixiv-backup

ensure_pixivpy3() {
    if /usr/bin/python3 -c "import pixivpy3" 2>/dev/null; then
        return 0
    fi

    echo "检测到缺少 pixivpy3，尝试自动安装..."
    if ! command -v pip3 >/dev/null 2>&1; then
        echo "错误: pip3 不可用，请先安装 python3-pip"
        return 1
    fi

    if ! pip3 install --no-cache-dir pixivpy3 >/tmp/pixiv-backup-pip.log 2>&1; then
        echo "错误: 自动安装 pixivpy3 失败，请检查网络和证书"
        tail -20 /tmp/pixiv-backup-pip.log 2>/dev/null
        return 1
    fi

    /usr/bin/python3 -c "import pixivpy3" 2>/dev/null
}

start_service() {
    if ! ensure_pixivpy3; then
        return 1
    fi

    # 检查配置
    if ! /usr/bin/python3 -c "import sys; sys.path.insert(0, '/usr/share/pixiv-backup'); from modules.config_manager import ConfigManager; c = ConfigManager(); print('Config valid' if c.validate_required() else 'Config invalid')" 2>&1 | grep -q "Config valid"; then
        echo "配置无效，请检查设置"
        return 1
    fi
    
    procd_open_instance
    procd_set_param command "$PROG" --daemon
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    # 需要读取 UCI 配置并写入输出目录，使用 root 运行避免权限导致崩溃
    procd_set_param pidfile /var/run/pixiv-backup.pid
    procd_set_param limits nofile="1024 2048"
    procd_close_instance
    
    echo "Pixiv备份服务已启动"
}

stop_service() {
    echo "停止Pixiv备份服务..."
    if ! /etc/init.d/pixiv-backup running >/dev/null 2>&1; then
        echo "服务未运行"
        return 0
    fi

    # 服务不存在时 procd_kill 会打印 ubus Not found，静默处理避免误报
    procd_kill "$PROG" >/dev/null 2>&1 || true
    echo "服务已停止"
}

restart() {
    stop
    sleep 1
    start
}

reload() {
    restart
}

EXTRA_COMMANDS="test status"
EXTRA_HELP="        test    测试配置和连接
        status  查看服务状态"

test() {
    echo "测试Pixiv备份配置..."
    
    # 检查Python依赖
    echo "检查Python依赖..."
    if ! ensure_pixivpy3; then
        echo "错误: pixivpy3 不可用"
        return 1
    fi
    
    if ! /usr/bin/python3 -c "import requests" 2>/dev/null; then
        echo "错误: requests未安装"
        return 1
    fi
    
    echo "Python依赖检查通过"
    
    # 测试配置
    echo "测试配置文件..."
    if ! /usr/bin/python3 -c "
import sys
sys.path.insert(0, '/usr/share/pixiv-backup')
from modules.config_manager import ConfigManager
try:
    config = ConfigManager()
    if config.validate_required():
        print('配置检查通过')
        print(f'用户ID: {config.get_user_id()}')
        print(f'输出目录: {config.get_output_dir()}')
    else:
        print('配置检查失败')
except Exception as e:
    print(f'配置错误: {e}')
    " 2>&1; then
        echo "配置测试失败"
        return 1
    fi
    
    # 测试API连接
    echo "测试Pixiv API连接..."
    if ! /usr/bin/python3 -c "
import sys
sys.path.insert(0, '/usr/share/pixiv-backup')
from modules.config_manager import ConfigManager
from modules.auth_manager import AuthManager
try:
    config = ConfigManager()
    auth = AuthManager(config)
    result = auth.test_connection()
    if result['success']:
        print('API连接测试通过')
        if 'user_name' in result:
            print(f'用户: {result[\"user_name\"]} (@{result[\"account\"]})')
    else:
        print(f'API连接失败: {result[\"error\"]}')
except Exception as e:
    print(f'连接测试错误: {e}')
    " 2>&1; then
        echo "API连接测试失败"
        return 1
    fi
    
    echo "所有测试通过！"
    return 0
}

status() {
    local pid
    pid=$(pgrep -f "pixiv-backup")
    
    if [ -n "$pid" ]; then
        echo "服务状态: 运行中 (PID: $pid)"
        
        # 显示运行信息
        if [ -f "/var/run/pixiv-backup.pid" ]; then
            echo "PID文件: /var/run/pixiv-backup.pid"
        fi
        
        # 显示日志文件（输出目录/data/logs 下最新日志）
        local log_file
        log_file=$(/usr/bin/python3 -c "
import sys
from pathlib import Path
sys.path.insert(0, '/usr/share/pixiv-backup')
from modules.config_manager import ConfigManager
config = ConfigManager()
log_dir = config.get_log_dir()
if not log_dir.exists():
    print('')
else:
    files = sorted(log_dir.glob('pixiv-backup-*.log'), key=lambda p: p.stat().st_mtime, reverse=True)
    print(files[0] if files else '')
" 2>/dev/null)
        if [ -f "$log_file" ]; then
            echo "日志文件: $log_file"
            echo "最近日志:"
            tail -5 "$log_file" 2>/dev/null || echo "无法读取日志"
        fi
    else
        echo "服务状态: 已停止"
    fi
    
    # 显示配置状态
    echo ""
    echo "配置状态:"
    /usr/bin/python3 -c "
import sys
sys.path.insert(0, '/usr/share/pixiv-backup')
from modules.config_manager import ConfigManager
try:
    config = ConfigManager()
    print(f'用户ID: {config.get_user_id() or \"未设置\"}')
    print(f'输出目录: {config.get_output_dir()}')
    print(f'下载模式: {config.get_download_mode()}')
    print(f'服务启用: {\"是\" if config.get(config.main_section, \"enabled\") == \"1\" else \"否\"}')
except Exception as e:
    print(f'读取配置失败: {e}')
    " 2>&1
    
    return 0
}
